# configurator/configs/looksee.upstart
start on static-network-up

respawn

chdir /opt/shmoocon_2014_talk

pre-start script
  echo "Updating looksee checkout..."
  git reset --hard HEAD
  git pull
end script

exec /usr/bin/python looksee/oversee.py

post-stop script
  echo "Sleeping 5 after stopping..."
  sleep 5
end script
