# configurator/configs/screenshot.conf
start on runlevel [2345]
stop on runlevel [!2345]

chdir /root/looksee

exec tasam 10 looksee:ScreenshotWorker
